# sample-configs/gochan-sqlite3.service
[Unit]
Description=gochan Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/gochan

[Install]
WantedBy=multi-user.target