# gochan.service
[Unit]
Description=gochan Daemon
After=network.target

[Service]
ExecStart=/home/vagrant/gochan/gochan

[Install]
WantedBy=multi-user.target